# systemd/taskflow-api.service
[Unit]
Description=Mission Control API Server
After=network.target postgresql.service

[Service]
Type=simple
User=openclaw
WorkingDirectory=/opt/mission-control/api
Environment="PATH=/opt/mission-control/venv/bin"
EnvironmentFile=/opt/mission-control/.env
ExecStart=/opt/mission-control/venv/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
